# departure-board.service
[Unit]
Description=LED Matrix Departure Board
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/departure-board
ExecStart=/usr/bin/env python3 /home/pi/departure-board/matrix_departure_board.py --stop "Basel, Aeschenplatz" --limit 4 --brightness 40 --gpio-mapping adafruit-hat
Restart=on-failure
RestartSec=5
# Nice value to reduce priority a bit
Nice=5

# Capability drop (not strictly required but safer); rgbmatrix needs GPIO (root usually). If using root, adjust.
# If library requires root, comment User/Group above and run as root.

[Install]
WantedBy=multi-user.target
