# usb-rescan.service
[Unit]
Description=Re-enumerate USB bus to recover devices missed at cold boot
# Run before any service that depends on USB peripherals (e.g. the USB speaker
# used by departure-board.service). udev should have populated the platform
# drivers by the time multi-user services start.
DefaultDependencies=yes

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/usb-rescan.sh
RemainAfterExit=yes
TimeoutStartSec=15

[Install]
WantedBy=multi-user.target
